# reporting/crontab
0 0 2 * * * /usr/local/bin/python /app/run_report_job.py
