# reporting/crontab
# Daily reporting tick (extend to call tenant report endpoints when available)
0 2 * * * python /app/run_report_job.py
